# .config/systemd/user/hyprpaper-switch.service
[Unit]
Description=Change wallpaper
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/hyprpaper_random
