# .config/systemd/user/swaybg-switch.timer
[Unit]
Description=Timer For swaybg random wallpaper

[Timer]
OnActiveSec=0
OnUnitActiveSec=15s
Persistent=true

[Install]
WantedBy=timers.target
